# app-vaultwarden/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio

# =============================================================================
# Start nginx.
# =============================================================================

bashio::net.wait_for 80
bashio::log.info "Starting nginx..."

exec nginx
