# vaultwarden/rootfs/etc/s6-overlay/s6-rc.d/init-nginx/run
#!/command/with-contenv bashio

# =============================================================================
# Configure nginx.
# =============================================================================

bashio::config.require.ssl

if bashio::config.true 'ssl'; then
    declare certfile=$(bashio::config 'certfile')
    declare keyfile=$(bashio::config 'keyfile')

    cp /etc/nginx/servers/direct-ssl.disabled /etc/nginx/servers/direct.conf

    sed -i "s#%%certfile%%#${certfile}#g" /etc/nginx/servers/direct.conf
    sed -i "s#%%keyfile%%#${keyfile}#g" /etc/nginx/servers/direct.conf
else
    cp /etc/nginx/servers/direct.disabled /etc/nginx/servers/direct.conf
fi
